# conf.d/blog.conf
[program:nodejs]
command=/usr/bin/node index.js
autostart=true
autorestart=true
stderr_logfile=/var/log/nodeapp_supervisor.err.log
stdout_logfile=/var/log/nodeapp_supervisor.out.log